# dot_config/systemd/user/flameshot.service
[Unit]
Description=Flameshot

[Service]
Type=simple
ExecStart=/usr/bin/flameshot

[Install]
WantedBy=default.target
