# dot_config/systemd/user/foot.service
[Unit]
Description=Foot Server

[Service]
Type=simple
ExecStart=/usr/bin/foot --server

[Install]
WantedBy=default.target
