# dot_config/systemd/user/ibus-daemon.service
[Unit]
Description=Ibus daemon

[Service]
Type=simple
ExecStart=/usr/bin/ibus-daemon -rx

[Install]
WantedBy=default.target
